# scripts/dashboard-update.service
[Unit]
Description=Dashboard Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-dashboard.sh
